# files/etc/init.d/blue-merle-hostname
#!/bin/sh /etc/rc.common

. /lib/blue-merle/functions.sh

START=10
STOP=99

start(){
    RANDOMIZE_HOSTNAME
}
stop(){
    echo "Stopping Blue-Merle Hostname Randomize Service"
}
